# systemd/dev-proxy-agent.service.in
[Unit]
Description=dev-proxy on-demand trigger agent
After=dev-proxy.service
Wants=dev-proxy.service

[Service]
Type=simple
ExecStart=@NODE_BIN@ @DEV_PROXY_ROOT@/trigger-agent.mjs
WorkingDirectory=@DEV_PROXY_ROOT@
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
